# cep2016.supervisor.ini
[program:cep2016]
command=/home/appmaster/cep2016web/env/bin/python /home/appmaster/cep2016web/env/bin/gunicorn web_cep_2016.wsgi
directory=/home/appmaster/cep2016web
user=appmaster
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cep2016.log
redirect_stderr=true
